# bash/.crontab
SHELL=/bin/bash
0 * * * * /Users/james.tomasino/bin/tmux-session save > /dev/null
