# bash/.crontab
0 * * * * . ~/.environment; bash /Users/james.tomasino/bin/tmux-session save > /dev/null
