# bash/.crontab
0 * * * * ~/bin/tmux-session save > /dev/null
